# conf/php56-fpm.service
[Unit]
Description=php56-fpm
After=network.target

[Service]
Type=forking
User=www
ExecStart="/usr/local/php56/sbin/php-fpm"
ExecStop="/bin/pkill php-fpm"
PrivateTmp=true
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
